# Philips Hue Bridge/Hardware/Firmware/_mtd11.bin.extracted/squashfs-root-0/etc/uci-defaults/z041_remove_updated_dummy_components
#!/bin/sh

rm -f /home/updated/componentinfo.json

exit 0