# Philips Hue Bridge/Hardware/Firmware/_mtd11.bin.extracted/squashfs-root-0/etc/uci-defaults/z047_migrate-firewall-remove-obsolete-eventing-rules
#!/bin/sh

if uci -q get firewall.eventing; then
	uci -q batch <<-EOF > /dev/null
		del firewall.eventing
		commit firewall
		EOF
fi

if uci -q get firewall.eventing_ipv6; then
	uci -q batch <<-EOF > /dev/null
		del firewall.eventing_ipv6
		commit firewall
		EOF
fi
