# Philips Hue Bridge/Hardware/Firmware/_mtd11.bin.extracted/squashfs-root/usr/lib/opkg/info/hk_hap.list
/etc/uci-defaults/z033_migrate-hk_hap-logging-interface-port-removal
/etc/init.d/hap
/etc/uci-defaults/z014_migrate-hk_hap-setup_from_uboot
/etc/uci-defaults/z031_migrate-hk_hap-product_data-to-file
/etc/uci-defaults/z024_migrate-hk_hap-logging-interface-address
/etc/uci-defaults/z016_migrate-hk_hap-hue_service
/etc/clip-resources/hap.json
/etc/uci-defaults/z005_migrate-hk_hap-add_device_id_file
/etc/uci-defaults/z011_migrate-hk_hap-upgrade
/etc/uci-defaults/z009_migrate-hk_hap-upgrade
/etc/uci-defaults/z033_migrate-hk_hap-natural_lighting_storage-from-file
/etc/uci-defaults/z046_migrate-hk_hap-homekit_fix_corrupted_srp_verifier
/etc/uci-defaults/z035_migrate-hk_hap-homekit_srp_verifier_workaround
/etc/uci-defaults/z032_migrate-hk_hap-default-product-data-to-file
/home/homekit/eui64_device_ids.tgz
/etc/uci-defaults/z034_migrate-hk_hap-homekit_resource_uuid-from-file
/etc/uci-defaults/z026_migrate-hk_hap-generated-srp-verifier-from-file
/etc/uci-defaults/z020_migrate-hk_hap-srv-dir-to-tmpfs
/etc/config/hk_hap
