# Philips Hue Bridge/Firmware/FW/_mtd9.bin.extracted/squashfs-root/usr/lib/opkg/info/hue-clip-daemon.list
/etc/clipd/POST_resource-matter_fabric.json
/etc/clipd/PUT_resource-button-id.json
/etc/clipd/PUT_resource-scene-id.json
/etc/clipd/POST_resource-behavior_instance.json
/etc/clipd/POST_resource-scene.json
/etc/clipd/PUT_resource-grouped_light-id.json
/etc/clipd/PUT_resource-clip-id.json
/etc/clipd/PUT_resource-speaker-id.json
/etc/clipd/PUT_resource-matter-id.json
/etc/clipd/PUT_resource-zgp_connectivity-id.json
/etc/clipd/PUT_resource-light-id.json
/etc/clipd/PUT_resource-room-id.json
/etc/clipd/PUT_resource-matter_fabric-id.json
/etc/clipd/PUT_resource-device_power-id.json
/etc/clipd/PUT_resource-behavior_instance-id.json
/etc/clipd/PUT_resource-wifi_connectivity-id.json
/etc/clipd/POST_resource-geofence_client.json
/etc/uci-defaults/z022_migrate-firewall-enable-eventing
/etc/clipd/PUT_resource-geofence_client-id.json
/etc/clipd/PUT_resource-tamper-id.json
/etc/clipd/POST_resource-service_group.json
/etc/clipd/PUT_resource-geolocation-id.json
/etc/clipd/PUT_resource-grouped_motion-id.json
/etc/clipd/POST_resource-room.json
/etc/clipd/PUT_resource-device-id.json
/etc/clipd/PUT_resource-security_area_motion-id.json
/etc/clip-resources/clip.json
/etc/clipd/PUT_resource-contact-id.json
/etc/clipd/PUT_resource-service_group-id.json
/etc/clipd/PUT_resource-camera_motion-id.json
/etc/clipd/PUT_resource-relative_rotary-id.json
/etc/clipd/PUT_resource-grouped_light_level-id.json
/etc/clipd/PUT_resource-bridge-id.json
/etc/init.d/clipd
/etc/clipd/PUT_resource-bell_button-id.json
/etc/clipd/PUT_resource-convenience_area_motion-id.json
/etc/clipd/PUT_resource-temperature-id.json
/etc/clipd/PUT_resource-zone-id.json
/etc/clipd/PUT_resource-zigbee_connectivity-id.json
/etc/clipd/PUT_resource-entertainment_configuration-id.json
/etc/clipd/PUT_resource-light_level-id.json
/etc/clipd/PUT_resource-zigbee_device_discovery-id.json
/etc/clipd/POST_resource-smart_scene.json
/etc/clipd/PUT_resource-motion_area_candidate-id.json
/etc/clipd/PUT_resource-entertainment-id.json
/etc/clipd/PUT_resource-motion-id.json
/etc/clipd/PUT_resource-smart_scene-id.json
/etc/clipd/PUT_resource-device_software_update-id.json
/etc/clipd/clip-api.schema.json
/etc/clipd/POST_resource-zone.json
/etc/clipd/PUT_resource-homekit-id.json
/etc/clipd/POST_resource-motion_area_configuration.json
/etc/clipd/PUT_resource-motion_area_configuration-id.json
/etc/clipd/POST_resource-entertainment_configuration.json
