# Philips Hue Bridge/Firmware/FW/_mtd9.bin.extracted/squashfs-root/usr/lib/opkg/info/platform-websocketcd.list
/etc/init.d/websocketcd
/usr/share/websocketcd_precondition_utils.sh
/etc/ca-certificates/ca.ecc.cert-and-crls.pem
/usr/bin/websocketcd
