# Philips Hue Bridge/Firmware/FW/_mtd11.bin.extracted/squashfs-root/etc/uci-defaults/z033_migrate-hk_hap-logging-interface-port-removal
#!/bin/sh

HK_LOGGING_INTERFACE_ADDRESS="localhost"

uci -q batch <<-EOF >/dev/null
        set hk_hap.hk_hap.logging_interface_address=${HK_LOGGING_INTERFACE_ADDRESS}
        commit hk_hap
EOF

exit 0
