# Philips Hue Bridge/Firmware/FW/_mtd11.bin.extracted/squashfs-root/etc/uci-defaults/z044_remove_existing_groups_database
#!/bin/sh

rm /home/groups/*.lmdb

exit 0
