# Philips Hue Bridge/Hardware/Firmware/_mtd9.bin.extracted/squashfs-root/etc/uci-defaults/z007_migrate-network-change_hostname
#!/bin/sh

[ "$(uci -q get network.lan.hostname)" != "Bridge-platform" ] && exit 0

uci -q batch <<-EOF >/dev/null
        set network.lan.hostname=Philips-hue
        commit network
EOF

exit 0
