# Philips Hue Bridge/Firmware/FW/_mtd11.bin.extracted/squashfs-root/usr/lib/opkg/info/iot-connectivity.list
/usr/bin/iot_connectivity/reprovisioning/__init__.py
/usr/bin/iot_connectivity/utils/__init__.py
/usr/bin/iot_connectivity/monitor_provisioning/__init__.py
/etc/init.d/iot-connectivity
/etc/ssl/certs/tunnel_pass/mqtt.hue.signify.com.cn.client.key
/usr/bin/iot_connectivity/jwt/__init__.py
/usr/bin/iot_connectivity/broker_configuration/config_builder.py
/usr/bin/iot_connectivity/utils/b64.py
/usr/bin/iot_connectivity/broker_configuration/expected_schema.py
/usr/bin/iot_connectivity/diagnostics/connectivity_stats.py
/usr/bin/iot-connectivity
/usr/bin/iot_connectivity/mqtt_bridge_state/__init__.py
/usr/bin/iot_connectivity/jwt/libcrypto.py
/usr/bin/iot_connectivity/diagnostics/client.py
/usr/bin/iot_connectivity/jwt/sign.py
/usr/bin/iot_connectivity/diagnostics/events.py
/usr/bin/iot_connectivity/diagnostics/periodic_reports.py
/usr/bin/iot_connectivity/diagnostics/stats_counter/__init__.py
/usr/bin/iot_connectivity/utils/types.py
/etc/ssl/certs/google_iot/0d69c7e1.0
/usr/bin/iot_connectivity/utils/schema_validator.py
/usr/bin/iot_connectivity/echo_message/__init__.py
/etc/ssl/certs/google_iot/main.pem
/etc/ssl/certs/google_iot/98889839.0
/etc/ssl/certs/tunnel_pass/mqtt.hue.signify.com.cn.client.ecc.cert.pem
/usr/bin/iot_connectivity/iot_connectivity.py
/etc/ssl/certs/google_iot/484642b6.0
/usr/bin/iot_connectivity/broker_configuration/conf_template.py
/usr/bin/iot_connectivity/broker_configuration/__init__.py
/etc/ssl/certs/google_iot/b0e59380.0
/usr/bin/iot_connectivity/__init__.py
/usr/bin/iot_connectivity/settings.py
/usr/bin/iot_connectivity/diagnostics/__init__.py
/usr/bin/iot_connectivity/diagnostics/connectivity_diagnostics.py
/usr/bin/iot_connectivity/utils/mosquitto.py
/etc/ssl/certs/google_iot/backup.pem
/usr/bin/iot_connectivity/shared.py
