# Philips Hue Bridge/Firmware/FW/_mtd9.bin.extracted/squashfs-root-0/usr/lib/opkg/info/hue-stream.list
/etc/init.d/stream
/etc/clip-resources/stream.json
/etc/uci-defaults/z027_migrate-firewall-open_stream_port_ipv6
/etc/uci-defaults/z017_migrate-firewall-open_stream_port
