# Philips Hue Bridge/Firmware/FW/_mtd9.bin.extracted/squashfs-root-0/usr/lib/opkg/info/platform-factory-reset.list
/etc/init.d/factoryreset
/lib/preinit/75_factory_reset
/etc/uci-defaults/update-factory-reset-counter
/usr/sbin/factoryreset.sh
/usr/sbin/factoryreset_daemon
