# Philips Hue Bridge/Firmware/FW/_mtd11.bin.extracted/squashfs-root-0/usr/lib/opkg/info/platform-crash-handler.list
/etc/sysctl.d/setuid.conf
/usr/libexec/anonymize
/etc/ca-certificates/sentry.pem
/etc/init.d/core_handler
/usr/libexec/core_handler
/etc/swbranch
/usr/bin/core_handler
