# Philips Hue Bridge/Firmware/FW/_mtd9.bin.extracted/squashfs-root-0/etc/uci-defaults/z003_migrate-network-change_hostname
#!/bin/sh

[ "$(uci -q get network.lan.hostname)" != "Philips-hue" ] && exit 0

uci -q batch <<-EOF >/dev/null
        set network.lan.hostname=Bridge-platform
        commit network
EOF

exit 0
