# Philips Hue Bridge/Firmware/FW/_mtd9.bin.extracted/squashfs-root-0/etc/uci-defaults/z042_matter-reset
#!/bin/sh

MATTER_CONFIG_PATH="/home/matter"

rm -rf ${MATTER_CONFIG_PATH}

exit 0