# Philips Hue Bridge/Hardware/Firmware/_mtd9.bin.extracted/squashfs-root-0/etc/uci-defaults/restore-factory-keys
#!/bin/sh
# Copyright (c) 2019 Signify Holding

ssh-factory-key -a -s
exit $?
