# Philips Hue Bridge/Hardware/Firmware/_mtd11.bin.extracted/squashfs-root-0/etc/crontabs/root
0 3 * * * /usr/libexec/count_mem
0 * * * * /usr/libexec/dnsmasq_watchdog check
0 0 * * * /usr/libexec/dnsmasq_watchdog force
