# Philips Hue Bridge/Firmware/FW/_mtd11.bin.extracted/squashfs-root-0/etc/uci-defaults/z041_migrate-matter-add-dev-dac
#!/bin/sh

MATTER_DAC_STR="{\"dac\":{\"private_key\":\"6siFy0E6JaAGKqDgipCqhZMQkQAjO19tdX8ke25o\/Vs=\",\"public_key\":\"BEVSZu6kmpeysVc2sq4gEqZ6ZH0OT0W6Jr96pRQXCzTz1KlC6iTG0evb9146MrOCZSlK8iYGGx7bfNvU1RB4fdU=\",\"dac\":\"MIIB5TCCAYugAwIBAgIIaNJyJlp3Aq8wCgYIKoZIzj0EAwIwPTElMCMGA1UEAwwcTWF0dGVyIERldiBQQUkgMHgxMDBCIG5vIFBJRDEUMBIGCisGAQQBgqJ8AgEMBDEwMEIwIBcNMjIwMzA4MDAwMDAwWhgPOTk5OTEyMzEyMzU5NTlaMFAxIjAgBgNVBAMMGU1hdHRlciBEZXYgREFDIDB4MTAwQi8weDIxFDASBgorBgEEAYKifAIBDAQxMDBCMRQwEgYKKwYBBAGConwCAgwEMDAwMjBZMBMGByqGSM49AgEGCCqGSM49AwEHA0IABEVSZu6kmpeysVc2sq4gEqZ6ZH0OT0W6Jr96pRQXCzTz1KlC6iTG0evb9146MrOCZSlK8iYGGx7bfNvU1RB4fdWjYDBeMAwGA1UdEwEB\/wQCMAAwDgYDVR0PAQH\/BAQDAgeAMB0GA1UdDgQWBBToWD5VjbfmSnx35wP\/ufuqbg1qyjAfBgNVHSMEGDAWgBTBB94N+FMUg+9\/p2jxMBs7Lp5swTAKBggqhkjOPQQDAgNIADBFAiEAmfxMy8C3KYlZVFW0vKOJjgyBerq4mfh4SJVqrsFAVvYCIDmRWIY7TtAyBTULdjkOUynYX6ZL9qUhW7x+YT8eR3ul\",\"pai\":\"MIIBtDCCAVugAwIBAgIIf6vEVZ+AzjIwCgYIKoZIzj0EAwIwGjEYMBYGA1UEAwwPTWF0dGVyIFRlc3QgUEFBMCAXDTIyMDMwODAwMDAwMFoYDzk5OTkxMjMxMjM1OTU5WjA9MSUwIwYDVQQDDBxNYXR0ZXIgRGV2IFBBSSAweDEwMEIgbm8gUElEMRQwEgYKKwYBBAGConwCAQwEMTAwQjBZMBMGByqGSM49AgEGCCqGSM49AwEHA0IABPc2krHvLgiicOY4zkOaVwd7ONEu0tQ0kv25FJkZv4kCaI\/olGsCG4YNE8gj64yzMlhDHWWRgacwn8C9ksI3\/eWjZjBkMBIGA1UdEwEB\/wQIMAYBAf8CAQAwDgYDVR0PAQH\/BAQDAgEGMB0GA1UdDgQWBBTBB94N+FMUg+9\/p2jxMBs7Lp5swTAfBgNVHSMEGDAWgBR4XOcFuGuPTm\/Hk6pgy0PqaWiC1TAKBggqhkjOPQQDAgNHADBEAiBkHIyYjqQGEPhG1uEiXrE+OGrCcko6DYyTAv3iULBgUwIgIOvJGXx7hAGWb8JQRFLO6FRs9QJziUzJTLBN3k3UaWA=\"}}"
MATTER_DAC_FILE_PATH="/home/matter/matterd_dac.cfg"

[ ! -d "/home/matter/" ] && mkdir -p "/home/matter/"
# The operation looks not power-fail safe. but there is migrate script z042_matter-reset that will delete the file.
# So we don't care about the data in the file at all
echo ${MATTER_DAC_STR} > ${MATTER_DAC_FILE_PATH}

exit 0
