# Philips Hue Bridge/Firmware/Directories/usr/lib/opkg/info/messenger.list
/etc/init.d/messenger
/usr/bin/messenger
