# Philips Hue Bridge/Firmware/Directories/usr/lib/opkg/info/platform-secure-console.list
/bin/secure-console.sh
/etc/init.d/secure-console
/etc/uci-defaults/restore-factory-keys
/lib/functions/secure-console.sh
/usr/sbin/ssh-factory-key
