# Philips Hue Bridge/Hardware/Firmware/_mtd9.bin.extracted/squashfs-root/usr/lib/opkg/info/updated.list
/usr/bin/update/iot_monitor.py
/usr/bin/update/component_bridge.py
/usr/bin/update/updated.py
/usr/bin/update/bridge/platform_bsb002.py
/usr/bin/update/__init__.py
/usr/bin/update/component_zigbee.py
/usr/bin/update/observer_diagnostics.py
/usr/bin/update/util/diagnostics.py
/usr/bin/update/util/mqtt_proxy.py
/usr/bin/update/bridge/ubi.py
/usr/bin/update/observer_iot.py
/usr/bin/update/component.py
/usr/bin/update/bridge/platform.py
/usr/bin/update/util/misc.py
/usr/bin/update/reactor_component_helper.py
/usr/bin/update/bridge/uboot.py
/usr/bin/update/util/diagnostics_proxy.py
/usr/bin/update/util/eventing.py
/usr/bin/updated
/usr/bin/update/util/stream_wrapper_form_data.py
/usr/bin/update/frontend_iot.py
/usr/bin/update/util/stream_wrapper.py
/usr/bin/update/util/stream_wrapper_curl.py
/usr/bin/update/config.py
/usr/bin/update/bridge/subprocess.py
/usr/bin/update/bridge/bootslot.py
/usr/bin/update/scheduler.py
/usr/bin/update/frontend_mqtt.py
/usr/bin/update/util/persistence.py
/etc/uci-defaults/z041_remove_updated_dummy_components
/usr/bin/update/reactor_zigbee_helper.py
/usr/bin/update/util/stream_wrapper_base.py
/usr/bin/update/bridge/platform_hse.py
/usr/bin/update/frontend_ipbridge.py
/usr/bin/update/bridge/__init__.py
/usr/bin/update/util/bytebuffer.py
/usr/bin/update/util/curl/c_api.py
/usr/bin/update/reactor.py
/usr/bin/update/zigbee_fw_fetcher.py
/usr/bin/update/util/http_server.py
/usr/bin/update/util/curl/ffex.py
/usr/bin/update/reactor_report_helper.py
/usr/bin/update/frontend_zigbee.py
/usr/bin/update/util/calculate_report.py
/usr/bin/update/frontend_sideload.py
/etc/init.d/updated
/usr/bin/update/util/singleton.py
